# conf/cron/analyze_images.cron
*/15 * * * * root /root/scripts/kuvat/analyze_images.sh >> /tmp/analyze_images/stout.log 2>> /tmp/analyze_images/sterr.log
